# roles/wireguard/templates/wireguard-resolve-dns.service.j2
[Unit]
Description=Wireguard Re-Resolve DNS of endpoints

[Service]
Type=oneshot
ExecStart=/usr/local/bin/wireguard-reresolve-dns {{ wg_interface }}